# etc/profile.d/vx-aliases-obtener-resolucion-pantalla.sh
alias resolucion-pantalla='vx-obtener-resolucion-pantalla'
